# initramfs/content/etc/inittab
null::sysinit:/bin/ln -s /dev/event0 /dev/touchscreen0

null::sysinit:/bin/mkdir -m 766 /var/run
null::sysinit:/bin/mkdir -m 766 /var/log

null::sysinit:/bin/echo >> /var/log/messages
null::respawn:/sbin/syslogd -n -D -S
null::respawn:/sbin/klogd -n

null::sysinit:/sbin/ifconfig lo up

tty1::askfirst:-/bin/ash
tty2::askfirst:-/bin/ash
tty3::askfirst:-/bin/ash
tty4::askfirst:-/bin/ash

tty5::askfirst:-/bin/hush
tty6::askfirst:-/bin/hush
tty7::askfirst:-/bin/hush
tty8::respawn:/bin/tail -f /var/log/messages

# Stuff to do when restarting the init process
::restart:/sbin/busybox init

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot

null::shutdown:/bin/killall klogd
null::shutdown:/bin/killall syslogd

::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a

